# services/low-battery-shutdown.service
[Unit]
Description=Power off when battery is critically low

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/low-battery-shutdown
